# chezmoi/dot_config/systemd/user/catppuccin-mode-sync.service
[Unit]
Description=Sync external Catppuccin theme with DMS light/dark mode

[Service]
Type=oneshot
ExecStart=%h/.local/bin/catppuccin-mode-sync
