# chezmoi/dot_config/systemd/user/catppuccin-mode-sync.path
[Unit]
Description=Watch DMS light/dark mode changes

[Path]
PathChanged=%h/.local/state/DankMaterialShell/session.json
Unit=catppuccin-mode-sync.service

[Install]
WantedBy=default.target
